# root/etc/services.d/myscript/run
#!/usr/bin/with-contenv bash

umask 022

python /config/scripy.p
